# debian/zfs-zed.zed.upstart
description     "OpenZFS Event Daemon (zed)"

start on local-filesystems

expect daemon
#respawn

exec /usr/sbin/zed
